# overlay/etc/services.d/postgres/finish
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Show helpful log message.
#======================================================================================================================

pg_ctl stop -D ${POSTGRESQL_DATA}
bf-svc-finish -s ${PWD##*/} -t
