# overlay/etc/services.d/postgres/finish
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Show helpful log message.
#======================================================================================================================

s6-setuidgid postgres ${POSTGRESQL_LIB}/finish
bf-svc-finish -s ${PWD##*/} -t
