# overlay/etc/s6-overlay/s6-rc.d/postgresql/finish
#!/usr/bin/nu

use bf
use bf-postgresql
use bf-s6
bf env load -x postgresql

# Stop server and terminate container, unless we are restoring from a backup
def main [...args] {
    # check whether or not we are restoring from a backup
    if (bf-postgresql restore is_restoring) {
        bf write debug "A restore operation is running: do not stop PostgreSQL or terminate the container."
        exit 0
    }

    # stop server gracefully
    bf-postgresql ctl stop

    # terminate container
    bf-s6 svc finish --terminate
}
