# overlay/etc/s6-overlay/s6-rc.d/postgresql/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-postgresql run ; run preflight"


#======================================================================================================================
# Run as PostgreSQL user
#   -D  path to the data directory
#======================================================================================================================

pg postgres -D "${BF_PG_DATA}"
